# phosphor-health-monitor.service.in
[Unit]
Description=BMC health monitoring
Wants=xyz.openbmc_project.EntityManager.service
After=xyz.openbmc_project.EntityManager.service

[Service]
ExecStart=@bindir@/health-monitor
Type=dbus
Restart=always
BusName=xyz.openbmc_project.HealthMon
SyslogIdentifier=phosphor-health-monitor

[Install]
WantedBy=multi-user.target
