# servicewarning.service.in
[Unit]
Description=Service level warning Action Script
[Service]
ExecStart=/usr/bin/service_warning_action.sh %I 
Type=oneshot
 
[Install]
WantedBy=multi-user.target
