# ipc/phosphor-ipc-monitor.service.in
[Unit]
Description=BMC IPC monitoring
After=obmc-mapper.target

[Service]
ExecStart=@bindir@/ipc-monitor
Type=simple
Restart=always
SyslogIdentifier=phosphor-ipc-monitor

[Install]
WantedBy=multi-user.target